# .config/systemd/user/compton.service
[Unit]
Description=Compton daemon

[Service]
ExecStart=/usr/bin/compton -C

[Install]
WantedBy=default.target
