# .config/systemd/user/gsd-xsettings.service
[Unit]
Description=Gnome settings daemon

[Service]
ExecStart=/usr/lib/gsd-xsettings

[Install]
WantedBy=default.target
