# .config/systemd/user/nextcloud.service
[Unit]
Description=Nextcloud

[Service]
ExecStart=%h/bin/nextcloud

[Install]
WantedBy=default.target
