# .config/systemd/user/nextcloud.service
[Unit]
Description=Nextcloud

[Service]
ExecStart=/home/ivan/bin/nextcloud

[Install]
WantedBy=default.target
