# .config/systemd/user/picom.service
[Unit]
Description=Picom daemon

[Service]
ExecStart=/usr/bin/picom -C

[Install]
WantedBy=default.target
