# .config/systemd/user/telegram.service
[Unit]
Description=Run telegram

[Service]
Type=oneshot
ExecStart=/usr/bin/telegram-desktop -startintray

[Install]
WantedBy=default.target
