# .config/systemd/user/xorg-config.service
[Unit]
Description=Xorg stuff

[Service]
ExecStart=/home/ivan/.config/i3/xorg-config

[Install]
WantedBy=default.target
