# .config/shell/profile.d/20-i3
export TERMINAL=st
export DESKTOP_SESSION=gnome
